# systemd/herald.service
[Unit]
Description=Herald feed fetcher daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/go/src/github.com/matthewjhunter/herald
ExecStart=%h/go/src/github.com/matthewjhunter/herald/herald daemon --config %h/go/src/github.com/matthewjhunter/herald/config/config.yaml --interval 5m
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
